# services/mt5_extended/supervisord.conf
[supervisord]
nodaemon=true

[program:mt5]
command=wine /opt/terminal/terminal.exe
autorestart=true
priority=10

[program:rpyc]
command=python3 /opt/mt5_rpyc_server.py
autorestart=true
priority=20